# cron_service/crontab
SHELL=/bin/bash
BASH_ENV=/container.env
33 2 * * * cd /app && /usr/local/bin/python daily_membership_update.py && /usr/local/bin/python daily_risk_update.py > /proc/1/fd/1 2>&1
